# deployment/socialmedia-health.service
[Unit]
Description=Social Media Automation Health Check Service
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=socialmedia
Group=socialmedia
WorkingDirectory=/opt/socialmedia
Environment="PATH=/opt/socialmedia/venv/bin"
ExecStart=/opt/socialmedia/venv/bin/python -m automation_stack.health --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/socialmedia

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=socialmedia-health

[Install]
WantedBy=multi-user.target
